# fabdeploy/config_templates/supervisor/celerybeat.conf
[program:{{ supervisor_prefix }}celerybeat]
command={{ active_env_link }}/bin/python {{ active_django_link }}/manage.py celerybeat
    --loglevel={{ loglevel }}
    --settings={{ settings }}
directory={{ active_var_link }}
user={{ user }}
stdout_logfile={{ active_log_link }}/{{ supervisor_prefix }}celerybeat.log
stderr_logfile={{ active_log_link }}/{{ supervisor_prefix }}celerybeat.log
priority=200
