# fabdeploy/config_templates/supervisor/celeryd.conf
[program:{{ supervisor_prefix }}celeryd]
command={{ active_env_link }}/bin/python {{ active_django_link }}/manage.py celeryd
    --loglevel={{ loglevel }}
    --settings={{ settings }}
directory={{ active_var_link }}
user={{ user }}
stdout_logfile={{ active_log_link }}/{{ supervisor_prefix }}celeryd.log
stderr_logfile={{ active_log_link }}/{{ supervisor_prefix }}celeryd.log
stopwaitsecs=600
priority=200
