# fabdeploy/config_templates/supervisor/celeryd.conf
[program:{{ supervisor_prefix }}celeryd]
command={{ current_env_link }}/bin/python {{ current_django_link }}/manage.py celeryd
    --loglevel={{ loglevel }}
    --settings={{ settings }}
directory={{ current_var_link }}
user={{ user }}
stdout_logfile={{ current_log_link }}/{{ supervisor_prefix }}celeryd.log
stderr_logfile={{ current_log_link }}/{{ supervisor_prefix }}celeryd.log
stopwaitsecs=600
priority=200
