# fabdeploy/config_templates/supervisor/celerybeat.conf
[program:{{ supervisor_prefix }}celerybeat]
command={{ current_env_link }}/bin/python {{ current_django_link }}/manage.py celerybeat
    --loglevel={{ loglevel }}
    --settings={{ settings }}
directory={{ current_var_link }}
user={{ user }}
stdout_logfile={{ current_log_link }}/{{ supervisor_prefix }}celerybeat.log
stderr_logfile={{ current_log_link }}/{{ supervisor_prefix }}celerybeat.log
priority=200
