# fabdeploy/config_templates/supervisor/gunicorn.conf
[program:{{ supervisor_prefix }}gunicorn]
command={{ active_env_link }}/bin/gunicorn_django
    --config {{ active_etc_link }}/gunicorn.conf.py
    {{ active_django_link }}/settings.py
directory={{ active_django_link }}
user={{ user }}
stdout_logfile={{ active_log_link }}/{{ supervisor_prefix }}gunicorn.log
stderr_logfile={{ active_log_link }}/{{ supervisor_prefix }}gunicorn.log
priority=300
