# fabdeploy/config_templates/supervisor/gunicorn.conf
[program:{{ supervisor_prefix }}gunicorn]
command={{ current_env_link }}/bin/gunicorn_django
    --config {{ current_etc_link }}/gunicorn.conf.py
    {{ current_django_link }}/settings.py
directory={{ current_django_link }}
user={{ user }}
stdout_logfile={{ current_log_link }}/{{ supervisor_prefix }}gunicorn.log
stderr_logfile={{ current_log_link }}/{{ supervisor_prefix }}gunicorn.log
priority=300
