# fabdeploy/config_templates/supervisor/redis_db.conf
[program:{{ supervisor_prefix }}redis_db]
command=redis-server {{ active_etc_link }}/redis_db.conf
directory={{ active_var_link }}
user={{ user }}
stdout_logfile={{ active_log_link }}/{{ supervisor_prefix }}redis_db.log
stderr_logfile={{ active_log_link }}/{{ supervisor_prefix }}redis_db.log
stopwaitsecs=600
priority=100
